# precise/gandalf/templates/gandalf-webserver.conf
start on runlevel [2345]
stop on runlevel [!2345]

setuid git
setgid git

respawn
respawn limit 10 5

exec /usr/local/bin/gandalf-webserver
