# precise/gandalf/templates/git-daemon.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec git daemon --base-path=/var/repositories --syslog --export-all &
